# deploy/systemd/techpulse-fetch.service
[Unit]
Description=Fetch TechPulse news

[Service]
Type=oneshot
User=techpulse
WorkingDirectory=/srv/techpulse/app
ExecStart=/srv/techpulse/venv/bin/python manage.py fetch_news
